# templates/default/kagent.service.erb
[Unit]
Description = Kagent, monitors/controls Hops services
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.kagent.user %>

PIDFile = <%= node[:kagent][:pid_file] %>
ExecStart = <%= "#{node[:kagent][:base_dir]}/bin/start-agent.sh" %>
ExecStop = <%= "#{node[:kagent][:base_dir]}/bin/stop-agent.sh" %>
ExecReload = <%= "#{node[:kagent][:base_dir]}/bin/restart-agent.sh" %>


Type = forking

[Install]
WantedBy = multi-user.target
