# templates/default/kagent.service.erb
[Unit]
Description = Kagent, monitors/controls Hops services
After = syslog.target network.target remote-fs.target <%= @deps %>
StartLimitIntervalSec=80
StartLimitBurst=10

[Service]
User = <%= node['kagent']['user'] %>
Restart=on-failure
RestartSec=5
PIDFile = <%= node['kagent']['pid_file'] %>
# Prefixing a command with '-' means that you should still continue, even if it fails.
ExecStartPre=-/bin/chown <%= node['kagent']['user'] %> <%= node['kagent']['pid_file'] %>
ExecStart = <%= "#{node['kagent']['base_dir']}/bin/start-agent.sh" %>
ExecStop = <%= "#{node['kagent']['base_dir']}/bin/stop-agent.sh" %>
ExecReload = <%= "#{node['kagent']['base_dir']}/bin/restart-agent.sh" %>

Type = forking

[Install]
WantedBy = multi-user.target
